# src/deb/control/prerm
#!/bin/bash

echo pre-remove

systemctl stop shinyproxy.service
systemctl disable shinyproxy.service