# src/deb/control/prerm
#!/bin/bash

service shinyproxy stop
if [ -d /run/systemd/system ];
then
  systemctl disable shinyproxy.service
fi
